# deploy/systemd/life-os-pull.service
[Unit]
Description=Pull Life-OS data tree from origin (fast-forward only)
After=network-online.target
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
User=life
Group=life
WorkingDirectory=/home/life/data
ExecStart=/usr/bin/git pull --ff-only --autostash --quiet
# --autostash handles the race where a bot write happens between pull
# cycles: stash, ff-pull, pop. Without it, pull aborts on dirty tree
# and Telegram alerts fire. With it, the next push picks up the stashed
# changes as normal dirty state.
# Non-fatal: a divergent local commit will fail ff-only; the push timer
# resolves that on its next run. Don't mark the unit failed for it.
SuccessExitStatus=0 1

NoNewPrivileges=true
PrivateTmp=true
