# deploy/systemd/life-os-dashboard.service
[Unit]
Description=Life-OS dashboard (FastAPI behind Caddy)
After=network-online.target
Wants=network-online.target
OnFailure=notify-failure@%n.service

[Service]
Type=simple
User=life
Group=life
WorkingDirectory=/home/life/app
EnvironmentFile=/home/life/app/.env
ExecStart=/home/life/app/venv/bin/uvicorn dashboard.app:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/life/app /home/life/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
