# deploy/systemd/life-os-morning.timer
[Unit]
Description=Run the Life-OS morning briefing daily

[Timer]
# 06:00 server local time. If the system was down, run on next boot.
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
