# deploy/systemd/life-os-app-pull.service
[Unit]
Description=Auto-deploy life-os-app — git fetch + pull + restart services if behind
After=network-online.target
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
User=life
Group=life
WorkingDirectory=/home/life/app
ExecStart=/bin/bash /home/life/app/deploy/bin/app-pull-and-restart.sh
