# deploy/systemd/life-os-push.service
[Unit]
Description=Push Life-OS data tree changes back to origin
After=network-online.target
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
User=life
Group=life
WorkingDirectory=/home/life/data
# Delegated to a script so systemd's specifier parser doesn't trip over
# date's %F/%T format. See deploy/bin/push-data-tree.sh for the body.
ExecStart=/bin/bash /home/life/app/deploy/bin/push-data-tree.sh
SuccessExitStatus=0 1

NoNewPrivileges=true
PrivateTmp=true
