# deploy/systemd/life-os-pull.service
[Unit]
Description=Pull Life-OS data tree from origin (fast-forward only)
After=network-online.target
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
User=life
Group=life
WorkingDirectory=/home/life/data
ExecStart=/usr/bin/git pull --rebase --autostash --quiet
# --rebase + --autostash handles BOTH races at once:
#   - dirty tree (bot mid-write):  autostash → pull → pop
#   - divergent local commits     (push timer made commits since last
#     successful pull, while origin advanced):  rebase replays them on
#     top of origin's tip
# --ff-only was too strict — it refused both cases and triggered the
# OnFailure alert. After two real-world failures (2026-06-04 14:25 and
# 15:40) the wider tolerance is the right tradeoff for this single-
# tenant box.
SuccessExitStatus=0 1

NoNewPrivileges=true
PrivateTmp=true
