# deploy/systemd/life-os-push.timer
[Unit]
Description=Push Life-OS data tree changes back every few minutes

[Timer]
# Offset from the pull timer so they don't fight.
OnBootSec=4min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
